# s6-user/sv/dwl-srv/run
#!/bin/execlineb -P
importas -i HOME HOME
cd ${HOME}
fdmove -c 2 1
exec dwl
