# s6-user/sv/mate-polkit-srv/run
#!/bin/execlineb -P
fdmove -c 2 1
exec /usr/lib/mate-polkit/polkit-mate-authentication-agent-1
