# s6-user/sv/restore-brightness/up
#!/bin/execlineb -P
if -tX { /usr/bin/isdirempty /sys/class/backlight }
	exec /sbin/light -I
